# infra/systemd/auto-retrain.service
[Unit]
Description=Quant System Automated Walk-Forward Retraining
After=network.target data-refresh.service

[Service]
Type=oneshot
WorkingDirectory=/quant_system
ExecStart=/usr/bin/python3 -m alpha.retrain.cli --horizons 24 --include-4h --sighup
User=ubuntu
Environment=PYTHONPATH=/quant_system
StandardOutput=append:/quant_system/logs/retrain_cron.log
StandardError=append:/quant_system/logs/retrain_cron.log
TimeoutStartSec=3600

[Install]
WantedBy=multi-user.target
