# infra/systemd/bybit-alpha.service
[Unit]
Description=Bybit Alpha Strategy Runner (ETH only (BTC removed — walk-forward FAIL))
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/quant_system
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -m scripts.run_bybit_alpha --symbols ETHUSDT --ws
Restart=always
RestartSec=30
StandardOutput=append:/quant_system/logs/bybit_alpha.log
StandardError=append:/quant_system/logs/bybit_alpha.log

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
