# infra/systemd/bybit-mm.service
[Unit]
Description=Bybit Market Maker (ORDIUSDT 50x)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/quant_system
EnvironmentFile=/quant_system/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -m scripts.run_bybit_mm --symbol ORDIUSDT --leverage 50 --order-size 100.0 --max-inventory 5000 --loss-limit 500 --gamma 0.3 --min-spread 2.0 --max-spread 20.0 --mode adaptive
Restart=on-failure
RestartSec=10
RestartPreventExitStatus=73
MemoryMax=256M

[Install]
WantedBy=multi-user.target
