# infra/systemd/data-refresh.service
[Unit]
Description=Quant System Daily Data Refresh
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/quant_system
ExecStart=/usr/bin/python3 -m data.downloads.data_refresh
User=ubuntu
Environment=PYTHONPATH=/quant_system
StandardOutput=append:/quant_system/logs/data_refresh.log
StandardError=append:/quant_system/logs/data_refresh.log
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
