# infra/systemd/health-watchdog.service
[Unit]
Description=Trading health watchdog check
After=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/quant_system
EnvironmentFile=/quant_system/.env
ExecStart=/usr/bin/python3 -m monitoring.watchdog --restart
StandardOutput=append:/quant_system/logs/health_watchdog.log
StandardError=append:/quant_system/logs/health_watchdog.log
