# infra/systemd/ic-decay-monitor.service
[Unit]
Description=IC Decay Monitor — daily alpha decay detection
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/quant_system
ExecStart=/usr/bin/python3 -m monitoring.ic_decay_monitor --alert
User=ubuntu
Environment=PYTHONPATH=/quant_system
StandardOutput=append:/quant_system/logs/ic_decay_monitor.log
StandardError=append:/quant_system/logs/ic_decay_monitor.log
