# infra/systemd/quant-backup.timer
[Unit]
Description=Quant system daily backup (4am UTC, after retrain)

[Timer]
OnCalendar=*-*-* 04:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
